# buildroot-external/rootfs-overlay/etc/profile.d/env.sh
#!/bin/sh

# 设置bash的PS1的显示

export PS1="\[\e[1;36m\][\u@\h \w]\\$\[\e[0m\]"