# buildroot-external/board/cm4sen/rootfs-overlay/etc/NetworkManager/dispatcher.d/40-route
GW=$(ip route | grep default | grep wwan0 | awk ' {print $3}')
if [ -n "$GW" ]; then
    ip route del default via $GW dev wwan0
    ip route add default via $GW dev wwan0 metric 0
fi