# buildroot-external/board/cm4sen/rootfs-overlay/usr/lib/systemd/system/rtc.service
[Unit]
Description=RTC Load and Sync Service
Before=systemd-timedated.service

[Service]
Type=oneshot
ExecStart=/sbin/hwclock -s
RemainAfterExit=true
StandardOutput=journal
ExecStop=/sbin/hwclock -w

[Install]
WantedBy=multi-user.target
